# plugins/bzfscron/crontab
*	*	*	*	*	/say "testing"
10-20/2	*/2	*	*	*	/say "hello"
1,5,6,8	1-4,8-12	*	*	*	/say "hello2"
0-10,20-30/2	*	*	*	*	/say "hello3"